# laurel.install
post_install() {
    echo ""
    echo -e "\033[1;33m  LAUREL: Run the \033[1mlaurel\033[1;33m command in your terminal to finish setup.\033[0m"
    echo ""
}

post_upgrade() {
    echo ""
    echo "==> Laurel upgraded. Restart the service to pick up changes:"
    echo "       systemctl --user restart laurel-replay"
    echo ""
}

pre_remove() {
    echo "==> Stopping Laurel..."
    echo "    Run: systemctl --user disable --now laurel-replay"
    echo ""
}
