# laurel-replay.service
[Unit]
Description=Laurel replay buffer (GPU Screen Recorder)
After=graphical-session.target

[Service]
Type=forking
ExecStart=%h/.local/bin/laurel-replay start
ExecStop=%h/.local/bin/laurel-replay stop
Environment=LAUREL_CONFIG=%h/.local/share/laurel/config.sh
Environment=DISPLAY=:0
Restart=on-failure
RestartSec=5

[Install]
WantedBy=graphical-session.target
